# scripts/hotprevue.service
[Unit]
Description=Hotprevue Backend
After=network.target postgresql.service

[Service]
Type=simple
User=kjell
WorkingDirectory=/opt/hotprevue/backend
EnvironmentFile=/opt/hotprevue/backend/.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
