# worker/hotprevue-worker.service
[Unit]
Description=Hotprevue AI worker (CLIP + ansiktsdeteksjon)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kjell
WorkingDirectory=/opt/hotprevue/worker
EnvironmentFile=/etc/hotprevue-worker.env
ExecStart=/opt/tenketank/venv/bin/python worker.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
